# templates/app.service.j2
{#
TODO: reuse your unit from Lab 1
TODO: Path to jar should use variables {{ app_dir }}/{{ app_file }}
TODO: Environment variable for PORT should use value from variable {{ app_port }}
#}
[Unit]
Description=DevOps_Lab3_Arthur_Pinchuk
Requires=network.target
After=network.target syslog.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=java -jar {{ app_dir }}/{{ app_file }}
User=arthur_pinchuk
Restart=on-failure
RestartSec=10
Environment=PORT={{ app_port }}
